# systemd/deltastack-options-snapshot-qqq.service
[Unit]
Description=DeltaStack QQQ intraday options snapshot capture
After=network.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/apps/deltastack
EnvironmentFile=/home/ec2-user/apps/deltastack/.env
ExecStart=/home/ec2-user/apps/deltastack/venv/bin/python scripts/capture_options_snapshot.py --underlying QQQ
StandardOutput=journal
StandardError=journal
