# systemd/deltastack-backup.timer
[Unit]
Description=Run DeltaStack backup daily at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
