# systemd/deltastack-healthcheck.timer
[Unit]
Description=Run DeltaStack health check daily at 00:45 UTC

[Timer]
OnCalendar=*-*-* 00:45:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
