# systemd/deltastack-options-snapshot-qqq.timer
[Unit]
Description=Capture QQQ options snapshots every 5 min during market hours

[Timer]
OnCalendar=Mon..Fri *-*-* 14:30..21:00:00/300
Persistent=false

[Install]
WantedBy=timers.target
