# systemd/deltastack-agent-madmax.service
[Unit]
Description=DeltaStack Mad Max agent run
After=network.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/apps/deltastack
EnvironmentFile=/home/ec2-user/apps/deltastack/.env
ExecStart=/home/ec2-user/apps/deltastack/venv/bin/python -m deltastack.agent run --agent mad_max
StandardOutput=journal
StandardError=journal
